# test/incus/bpfrxd.service
[Unit]
Description=bpfrx eBPF Firewall Daemon
After=network-online.target frr.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/bpfrxd
Restart=on-failure
RestartSec=1
TimeoutStopSec=20
LimitMEMLOCK=infinity
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bpfrxd

[Install]
WantedBy=multi-user.target
