# test/incus/xpfd.service
[Unit]
Description=xpf eBPF Firewall Daemon
After=network-online.target frr.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/xpfd
Restart=on-failure
RestartSec=1
TimeoutStopSec=20
LimitMEMLOCK=infinity
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xpfd

# #712 Option A retry (after #740 fix): measured on the 6-core loss
# userspace lab with workers correctly pinned to CPUs 2-5 via
# CPUAffinity=2 3 4 5 — no metric moved by the #712 thresholds
# (rate ratio +2%, retrans +11%, per-flow CoV mean -0.6 pp, per-flow
# CoV max +1.7 pp; all within run-to-run noise). The directive is
# NOT shipped. See docs/cos-validation-notes.md §"CPU pinning retry
# post-#740" for the dated measurement, docs/712-cpu-pinning-recipe.md
# for the recipe (kept as design intent), and #739 for the kernel
# cmdline (isolcpus/nohz_full) follow-up which is now the next lever
# to try on this hardware.

[Install]
WantedBy=multi-user.target
