# test/incus/xpfd.service
[Unit]
Description=xpf eBPF Firewall Daemon
After=network-online.target frr.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/xpfd
Restart=on-failure
RestartSec=1
TimeoutStopSec=20
LimitMEMLOCK=infinity
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xpfd

# #712 Option A (CPUAffinity=) was measured on the 6-core loss userspace
# lab and did not move per-flow jitter, retransmits, or rate ratio on
# this hardware; no directive is added here. The cause is that the
# userspace-dp helper calls sched_setaffinity per worker (see
# pin_current_thread() in userspace-dp/src/afxdp/neighbor.rs) and that
# call runs AFTER systemd's CPUAffinity= is applied, so the hot-path
# workers land on CPUs 0..N-1 regardless of the unit-level mask. Only
# the Go daemon and the dp auxiliary threads honoured CPUAffinity=.
# See docs/712-cpu-pinning-recipe.md for the per-CPU-budget recipe
# (which is still useful on hosts where the worker count is ≤ the
# budget) and the measured-no-op finding that blocks it on 6-core.

[Install]
WantedBy=multi-user.target
